# systemd/openhands-symphony.service
[Unit]
Description=OpenHands Symphony GitHub issue orchestrator
After=network-online.target openhands-symphony-firewall.service openhands-canvas.service
Requires=openhands-symphony-firewall.service openhands-canvas.service
PartOf=openhands-symphony.target

[Service]
Type=simple
User=openhands-symphony
Group=openhands-symphony
SupplementaryGroups=openhands-agents
Environment=HOME=/var/lib/openhands-symphony
Environment=PATH=/opt/provider-clis/node_modules/.bin:/opt/openhands-acp/node_modules/.bin:/var/lib/openhands-symphony/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=DISABLE_AUTOUPDATER=1
Environment=SYMPHONY_CONFIG=/etc/openhands-symphony/config.toml
Environment=GH_CONFIG_DIR=/var/lib/openhands-symphony/github
ExecStart=/opt/openhands-symphony-tool/bin/openhands-symphony
Restart=on-failure
RestartSec=10
# The service uses one narrowly-scoped sudoers transition to the lower-authority
# credential-free validator. no_new_privs would disable that setuid transition.
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/openhands-symphony /var/lib/openhands-validator /var/log/openhands-symphony
UMask=0077

[Install]
WantedBy=openhands-symphony.target
